# dockmind-egpu-unbind.service
[Unit]
Description=Unbind DockMind eGPU drivers

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/dockmind-egpu-unbind
User=root
Group=root
